# package/remux/prerm
#!/usr/bin/env bash

echo "Disabling remux"
systemctl enable remux
systemctl start remux
